# systemd/paracoding-autosync.service
[Unit]
Description=fleet agent-workspace git autosync
[Service]
Type=oneshot
ExecStart=/opt/paracoding-mcp/autosync.sh
